# deploy/boss.service
[Unit]
Description=BOSS Mini-App Platform
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=rpi
WorkingDirectory=/opt/boss
ExecStart=/opt/boss/.venv/bin/python -m boss.main
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=graphical.target
