# deploy/boss.service
[Unit]
Description=BOSS Mini-App Platform
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=rpi
WorkingDirectory=/opt/boss
ExecStart=/opt/boss/.venv/bin/python -m boss.main
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/boss/secrets/secrets.env
SupplementaryGroups=gpio

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/boss/logs /opt/boss/secrets
MemoryMax=512M
CPUQuota=80%

[Install]
WantedBy=graphical.target
